# logstash.conf
[program:logstasah]
command=/opt/logstash/bin/logstash agent -f /etc/logstash/conf.d
user=logstash
redirect_stderr=true
directory=/var/lib/logstash
